# cron.d/adial-endpoints-refresh
# Refresh endpoints cache every 5 minutes
# This keeps the SIP/PJSIP endpoint list up to date without slowing down the web interface

*/5 * * * * root /usr/bin/php /var/www/html/adial/refresh-endpoints-cache.php >> /var/www/html/adial/logs/endpoints-refresh.log 2>&1
